# deploy/rc.local
#copiar essa linha no arquivo /etc/rc.local
/usr/local/enviar_sms_digivoice/scripts/montar_unidade.sh